# ci/ansible/roles/zenotech.slurm/files/munge.service
[Service]
ExecStart=
ExecStart=/usr/sbin/munged --syslog
